# pcg-study/設定ファイル/HWEB001/Route53_Aレコード自動書き換えスクリプト/ec2-route53.service
Decription = EC2 Dynamic IP Route 53 auto register

[Service]
ExecStart = /var/config/aws/ec2Setup.sh
Restart = always
Type = simple

[Install]
WantedBy = default.target
